# arch/gambit/files/etc/systemd/system/zfs-auto-snapshot-weekly.service
[Unit]
Description=ZFS weekly snapshot service

[Service]
ExecStart=/bin/zfs-auto-snapshot --skip-scrub --prefix=znap --label=weekly --keep=12 //
