# arch/common/files/etc/systemd/system/libvirtd.service.d/snapshot.conf
[Service]
ExecStartPost=-/usr/local/bin/virsh-dump-all.sh
ExecStartPost=-/usr/bin/zfs-auto-snapshot --skip-scrub --prefix=znap --label=libvirtd --keep=10 --recursive z/images
