# popos/common/files/etc/systemd/system/zfs-auto-snapshot-weekly.timer.d/override.conf
[Timer]
OnCalendar=
OnCalendar=Sun *-*-* 3:00
