# arch/common/files/etc/systemd/system/swap-ntfs.service
[Unit]
Description=Swap shared with Windows

[Service]
Type=notify
ExecStart=/usr/local/bin/swapon.sh
ExecStop=/usr/local/bin/swapoff.sh
RemainAfterExit=true
TimeoutStopSec=600

[Install]
WantedBy=multi-user.target
