# ubuntu/common-pop_os/files/etc/systemd/system/libvirtd.service.d/snapshot.conf
[Service]
ExecStartPre=-/usr/local/bin/zfs-auto-snapshot --skip-scrub --prefix=znap --label=libvirtd --keep=10 --recursive root/images