# arch/common/files/etc/systemd/system/zfs-auto-snapshot-monthly.service.d/override.conf
[Service]
ExecStart=
ExecStart=/usr/bin/zfs-auto-snapshot --skip-scrub --prefix=znap --label=monthly --keep=12 //
