# ubuntu/common/files/etc/systemd/system/update.service
[Unit]
Description=Auto-Update

[Service]
ExecStart=/usr/local/bin/update.sh
