# ubuntu/common/files/etc/systemd/system/zfs-replicate.service
[Unit]
Description=zfs replication

[Service]
Nice=19
IOSchedulingClass=idle
KillSignal=SIGINT
ExecStart=/usr/local/bin/sync-to-nas.sh
