# arch/strange/files/etc/systemd/system/zfs-auto-snapshot-daily.timer.d/override.conf
[Timer]
OnCalendar=
OnCalendar=*-*-* 3:00
