# arch/common/files/etc/systemd/system/zfs-auto-snapshot-frequent.service.d/override.conf
[Service]
ExecStart=
ExecStart=/bin/zfs-auto-snapshot --skip-scrub --prefix=znap --label=frequent --keep=32 //
