# arch/strange/files/etc/systemd/system/zfs-auto-snapshot-monthly.timer.d/override.conf
[Timer]
OnCalendar=
OnCalendar=Sun *-*-1..7 3:00:00
