# ubuntu/common-pop_os/files/etc/systemd/system/docker-prune.timer
[Unit]
Description=Prune Docker

[Timer]
OnCalendar=Weekly
RandomizedDelaySec=1h
Persistent=true

[Install]
WantedBy=timers.target
