# arch/vm-a/files/etc/systemd/system/zfs-replicate.service
[Unit]
Description=zfs replication

[Service]
Nice=19
IOSchedulingClass=idle
KillSignal=SIGINT
ExecStart=/root/opt/bin/vm-a/sync-to-beast.sh
