# ubuntu/common/files/etc/systemd/system/docker.service.d/snapshot.conf
[Service]
ExecStartPre=-/usr/local/bin/zfs-auto-snapshot --skip-scrub --prefix=znap --label=docker --keep=10 --recursive root/docker root/volumes