# ubuntu/common-min/files/etc/systemd/system/zfs-auto-snapshot-frequent.service
[Unit]
Description=ZFS frequent snapshot service

[Service]
ExecStart=/usr/local/bin/zfs-auto-snapshot --skip-scrub --prefix=znap --label=frequent --keep=8 //
