# popos/common/files/etc/systemd/system/zfs-auto-snapshot-monthly.service
[Unit]
Description=ZFS monthly snapshot service

[Service]
ExecStart=/usr/sbin/zfs-auto-snapshot --skip-scrub --prefix=znap --label=monthly --keep=12 //
