# arch/common/files/etc/systemd/system/zfs-auto-snapshot-hourly.service.d/override.conf
[Service]
ExecStart=
ExecStart=/usr/bin/zfs-auto-snapshot --skip-scrub --prefix=znap --label=hourly --keep=36 //
