# arch/common/files/etc/systemd/user/smbnetfs.service
[Unit]
Description=smbnetfs

[Service]
ExecStart=/usr/bin/smbnetfs %h/smb
ExecStop=/bin/fusermount -u %h/smb

[Install]
WantedBy=default.target
